# deploy/launchd/com.traderbot.stop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.traderbot.stop</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/macmini/projects/traderbot/deploy/launchd/stop-paper.sh</string>
    </array>

    <!-- Fire Friday at 17:00 local time. Weekday=5 is Friday in launchd. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/macmini/projects/traderbot/logs/launchd-stop.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/macmini/projects/traderbot/logs/launchd-stop.log</string>
</dict>
</plist>
